# biomine-ping.service
[Unit]
Description=Biomine Control Center

[Service]
ExecStart=/etc/biomine-scripts/ping.sh
